# agent/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; Virtual display :99
[program:xvfb]
command=Xvfb :99 -screen 0 1280x800x24 -ac
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

; Chromium on :99, CDP on port 9222
[program:chromium]
command=/usr/bin/chromium --no-sandbox --disable-gpu --disable-dev-shm-usage --disable-crash-reporter --window-size=1280,800 --window-position=0,0 --start-maximized --remote-debugging-port=9222 --remote-debugging-address=127.0.0.1 --password-store=basic --no-first-run --disable-search-engine-choice-screen --user-data-dir=/home/rcm/chrome
environment=DISPLAY=":99",HOME="/home/rcm"
autostart=true
autorestart=true
startretries=999
priority=10
user=rcm
stdout_logfile=/var/log/chromium.log
stderr_logfile=/var/log/chromium.err

; x11vnc — captures :99
[program:x11vnc]
command=x11vnc -display :99 -passwd %(ENV_VNC_PASSWORD)s -forever -shared -rfbport 5900
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

; mcp-proxy — exposes chrome-devtools-mcp over HTTP (SSE/StreamableHTTP)
[program:mcp-proxy]
command=mcp-proxy --port 3000 --host 127.0.0.1 --pass-environment -- npx -y chrome-devtools-mcp@latest --browserUrl http://127.0.0.1:9222
environment=HOME="/home/rcm"
autostart=true
autorestart=true
startretries=999
priority=30
user=rcm
stdout_logfile=/var/log/mcp-proxy.log
stderr_logfile=/var/log/mcp-proxy.err

; gateway — reverse proxy with token auth for MCP
[program:gateway]
command=/usr/local/bin/remote-chrome-mcp
environment=AUTH_TOKEN="%(ENV_AUTH_TOKEN)s",LISTEN_ADDR=":8080"
autostart=true
autorestart=true
startretries=999
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
